# roles/img-panda/templates/supervisord.conf.j2
[program:{{ app_name }}]
command=/usr/bin/gunicorn -b 0.0.0.0:{{ app_port }} {{ app_name }}:app
directory=/home/vagrant/{{ app_name }}-app
user=vagrant
autorestart=true
redirect_stderr=true
